# contrib/zerotier-systemd-manager.timer
[Unit]
Description=Update zerotier per-interface DNS settings

[Timer]
OnStartupSec=1min
OnUnitInactiveSec=1min

[Install]
WantedBy=timers.target
