# contrib/zerotier-systemd-manager.timer
[Unit]
Description=Update zerotier per-interface DNS settings

[Timer]
OnStartupSec=2min
OnUnitInactiveSec=1min

[Install]
WantedBy=timers.target
